# com.datashaman.gofer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.datashaman.gofer</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace these paths before installing -->
        <string>UV_BINARY</string>
        <string>run</string>
        <string>--project</string>
        <string>PROJECT_DIR</string>
        <string>gofer</string>
        <string>--config</string>
        <string>PROJECT_DIR/config.yaml</string>
        <string>--log-file</string>
        <string>LOG_DIR/gofer.log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PROJECT_DIR</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>LOG_DIR/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/stderr.log</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
