# com.datashaman.jira-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.datashaman.jira-agent</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the absolute path to your uv binary -->
        <string>/Users/YOU/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/YOU/Projects/datashaman/jira-agent</string>
        <string>jira-agent</string>
        <string>--config</string>
        <string>/Users/YOU/Projects/datashaman/jira-agent/config.yaml</string>
        <string>--log-file</string>
        <string>/Users/YOU/Library/Logs/jira-agent/jira-agent.log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/Projects/datashaman/jira-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/YOU/.local/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOU/Library/Logs/jira-agent/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/Library/Logs/jira-agent/stderr.log</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
